# scripts/sankhya_ajuda_sync.cron
# /etc/cron.d/sankhya_ajuda_sync
# Daily full sync of the Sankhya Help Center at 03:00 server time.
#
# flock holds a lock shared with the community sync (04:00). The two ETLs share
# the same Postgres + vLLM endpoint; the lock guarantees they never run at the
# same time even if this job overruns its hour. `-n` makes the OTHER job skip
# (and self-heal on its next daily run) rather than queue up behind a stuck one.
#
# Adjust the absolute path below to wherever you cloned the repo.
# Install with:
#   sudo cp scripts/sankhya_ajuda_sync.cron /etc/cron.d/sankhya_ajuda_sync
#   sudo chmod 644 /etc/cron.d/sankhya_ajuda_sync
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

0 3 * * * root flock -n /var/lock/sankhya_ajuda_etl.lock bash -c 'cd /path/to/sankhya-ajuda-mcp && .venv/bin/python -m sync.sync --full' >> /var/log/sankhya_ajuda_sync.log 2>&1
